# deploy/whatsapp-proxy.conf
[program:whatsapp-proxy]
command=/home/ubuntu/whatsapp-proxy/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8001 --workers 2
directory=/home/ubuntu/whatsapp-proxy
user=ubuntu
autostart=true
autorestart=true
startretries=3
stopwaitsecs=10
stderr_logfile=/var/log/whatsapp-proxy.err.log
stdout_logfile=/var/log/whatsapp-proxy.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=PYTHONUNBUFFERED="1"
